# root/etc/services.d/cron/run
#!/usr/bin/execlineb -P
# shellcheck shell=sh
cron -f -L 15
